# sv/log/run
#!/bin/sh

name="hubelmeter"
dir="${LOGDIR}/${name}"

mkdir -p "${dir}"
exec svlogd -tt "${dir}"
